# supervisor.conf
[supervisord]
nodaemon = true

[program:scrapyd]
command=/bin/bash -c 'scrapyd' -u scrapy -g scrapy -l /var/log/scrapyd/scrapyd.log
stdout_events_enabled = true
stderr_events_enabled = true
